# private_dot_config/systemd/user/notify-hydration.service
[Unit]
Description=Personal — hydration nudge

[Service]
Type=oneshot
ExecStart=%h/.local/share/chezmoi/system_scripts/notify/send-personal.sh hydration
NoNewPrivileges=true
