# private_dot_config/systemd/user/notify-host-monitors.timer
[Unit]
Description=Run host monitors (disk/thermal/network) every 5 minutes
Documentation=file:///home/shantanu/.local/share/chezmoi/docs/system/Host Monitors.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=notify-host-monitors.service

[Install]
WantedBy=timers.target
