# private_dot_config/systemd/user/notify-posture.service
[Unit]
Description=Personal — posture/stretch nudge

[Service]
Type=oneshot
ExecStart=%h/.local/share/chezmoi/system_scripts/notify/send-personal.sh stretch
NoNewPrivileges=true
